# docker/php/cron.d/symfony.development
# ==========
# MAINTENANCE
# ==========

# Logs the last cron execution
*  *  *  *  *  www-data  date > /var/www/web/export/last-cron
